# pi-gen-sources/00-teslausb-tweaks/files/backup.service
[Unit]
Description=run force sync

[Service]
Type=simple
ExecStart=/root/bin/force_sync.sh

[Install]
WantedBy=default.target
